# root/etc/service/xinetd/run
#!/bin/sh
# Start glftpd
exec /usr/sbin/xinetd -pidfile /run/xinetd.pid -stayalive -inetd_compat -inetd_ipv6 -dontfork

